# scripts/com.prevoir.poll-jira.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS launchd schedule for poll-jira.sh
  Runs every 60 minutes; compatible with Power Nap when plugged in.

  SETUP:
  1. Replace /Users/YOUR_USERNAME below with your actual home directory path
  2. Copy this file to ~/Library/LaunchAgents/
  3. Run: launchctl load ~/Library/LaunchAgents/com.prevoir.poll-jira.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.prevoir.poll-jira</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <!-- UPDATE: replace /Users/YOUR_USERNAME with your actual home path -->
        <string>/Users/YOUR_USERNAME/Documents/Prevoir/Scripts/poll-jira.sh</string>
    </array>

    <!-- Run every 60 minutes; works with Power Nap when plugged in -->
    <key>StartInterval</key>
    <integer>3600</integer>

    <!-- UPDATE: replace /Users/YOUR_USERNAME with your actual home path -->
    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/Documents/Prevoir/Scripts/poll-jira.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/Documents/Prevoir/Scripts/poll-jira-error.log</string>

    <!-- Only run when the user is logged in -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
